# configuration/arl.service
[Unit]
Description=ARL

[Service]
ExecStart=/opt/arl/arl -c /opt/arl/arl.config

[Install]
WantedBy=multi-user.target
